# Task05/85-Task05.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", RUN+="/home/parham/bin/Task01.sh"
